# martin-luther-so/fs/etc/crontabs/root
23 * * * *	/usr/sbin/owm.lua
*/4 * * * * /usr/bin/auto_ipv6_node
*/5 * * * *	killall -HUP dnsmasq
*/5 * * * *	/usr/sbin/ff_olsr_watchdog
* * * * *	/usr/sbin/ffwatchd running || /etc/init.d/freifunk-watchdog restart
