# tberg-core/fs/etc/crontabs/root
*/5 * * * *	killall -HUP dnsmasq
*/5 * * * *	/usr/sbin/ff_olsr_watchdog
# reboot if country could not be set
*/10 * * * *       (iw reg get | grep -q DFS-UNSET) && reboot
