# nkarcaden-fhain/fs/etc/crontabs/root
*/4 * * * * /usr/bin/auto-ipv6
*/5 * * * *     /usr/sbin/ff_luci_resolve
* * * * *	/usr/sbin/ff_olsr_test_gw.sh
41 0 * * * /sbin/heartbeat
37 * * * *	/usr/sbin/owm.lua
*/5 * * * *	killall -HUP dnsmasq
*/5 * * * *	/usr/sbin/ff_olsr_watchdog
* * * * *	/usr/sbin/ffwatchd running || /etc/init.d/freifunk-watchdog restart
*/5 * * * * 	/usr/sbin/luci-splash sync
