# segen-router/fs/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.


( sleep 60 ; /etc/init.d/dnsmasq restart ) &
( sleep 62 ; /etc/init.d/olsrd restart ) &

exit 0
