# elok-lascafe/fs/etc/crontabs/root
*/5 * * * * ( for file in /etc/alfred/* ; do [ -x $file ] && $file ; done )
