# emma-wsw-m2/fs/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

batctl if add $(uci -p /var/state get network.lan.ifname)
batctl if add $(uci -p /var/state get network.wireless0.ifname)

exit 0
